# docker/bitloudd.sv.conf
[program:apache]
command=/usr/sbin/apachectl -D FOREGROUND
autostart=true
autorestart=true
redirect_stdout=true
redirect_stderr=true
startsecs=30
